# daemon/harbour-lagoon-daemon.service
[Unit]
Description=Lagoon Background Daemon
After=network.target dbus.service

[Service]
Type=dbus
BusName=org.harbour.lagoon
ExecStart=/usr/bin/harbour-lagoon-daemon
Restart=on-failure
RestartSec=10

# Resource limits
MemoryLimit=100M
CPUQuota=20%

# Security
NoNewPrivileges=true
PrivateTmp=true

# Environment
Environment="QT_QPA_PLATFORM=offscreen"
Environment="QT_LOGGING_RULES=*.debug=false"

[Install]
WantedBy=default.target
